# type/__wikijs/files/wikijs-openrc
#!/sbin/openrc-run

command='/usr/bin/node'
command_args='server'
command_background=true
command_user='wikijs'
description="Run wiki.js"
directory='/var/wiki'
error_log=/var/log/"$RC_SVCNAME".err
output_log=/var/log/"$RC_SVCNAME".log
pidfile="/run/$RC_SVCNAME.pid"
